# autofs/files/autofs.init
#!/bin/sh /etc/rc.common
#
# Author: Zhifeng Gu <guzhifeng1979@hotmail.com>
# Date: 2012/03
# License: GPL

START=21
STOP=98
AUTOMOUNT_PIDFILE=/var/run/autofs.pid

start() {
    /usr/sbin/automount -p ${AUTOMOUNT_PIDFILE}
}

stop() {
    kill -TERM $(cat ${AUTOMOUNT_PIDFILE})
}

reload() {
    kill -HUP $(cat $AUTOMOUNT_PIDFILE)
}
